# include_bin/services/check-ip-change.service
[Unit]
Description=Check for IP change on boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check_ip_change.sh

[Install]
WantedBy=multi-user.target
