# include_bin/services/full-backup.service
[Unit]
Description=Full Pi Backup Image

[Service]
Type=simple
ExecStart=/usr/local/bin/full_backup_image.sh
